# debian/approx.cron.weekly
#!/bin/sh

# Garbage collect the approx(8) cache

/usr/sbin/gc_approx --quiet
