# debian/approx.cron.daily
#!/bin/sh

# Update the approx(8) cache

if [ -x /usr/sbin/update_approx ]; then
    /usr/sbin/update_approx --quiet
fi
