# debian/approx.cron.weekly
#!/bin/sh

# Garbage collect the approx(8) cache

gc_approx --quiet
