# debian/approx.cron.weekly
#!/bin/sh

# Garbage collect the approx(8) cache

if [ -x /usr/sbin/gc_approx ]; then
    /usr/sbin/gc_approx --quiet
fi
