# debian/approx.cron.weekly
#!/bin/sh

# Garbage collect the approx(8) cache

# Run as the approx user in case we need to download replacement files

su approx -c "gc_approx --quiet"

# Remove any empty directories
# The trailing /. makes it work when /var/cache/approx is a symlink

find /var/cache/approx/. -type d -empty | xargs rmdir
