# debian/approx.cron.daily
#!/bin/sh

# Update the approx(8) cache

if [ -x /usr/sbin/approx-update ]; then
    su --shell /bin/sh --command "/usr/sbin/approx-update --quiet" approx
fi
