# debian/approx.cron.weekly
#!/bin/sh

# Garbage collect the approx(8) cache

if [ -x /usr/sbin/approx-gc ]; then
    su --shell /bin/sh --command "/usr/sbin/approx-gc --quiet" approx
fi
