# debian/approx.cron.weekly
#!/bin/sh

# Garbage collect the approx(8) cache

if [ -x /usr/sbin/approx-gc ]; then
    nice -n19 ionice -c3 /usr/sbin/approx-gc --quiet
fi
